# pi/sky-sentry.service
[Unit]
Description=Sky Sentry: flash Arduino + run timelapse capture
# Make sure the SSD is mounted before we try to write frames there.
RequiresMountsFor=/mnt/ssd
# Wait for the network only if you actually need it (e.g. NTP-correct timestamps).
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Edit these three to match your setup:
User=pi
Group=pi
WorkingDirectory=/home/pi/Timelapse-Lightning-Setup

# Optional: pass flags through to timelapse.py.
# Examples:
#   Environment=TIMELAPSE_ARGS=--interval 30
#   Environment=ARDUINO_PORT=/dev/ttyACM0
#   Environment=SKIP_FLASH=1
Environment=TIMELAPSE_ARGS=

ExecStart=/home/pi/Timelapse-Lightning-Setup/pi/sky-sentry.sh

# Graceful Ctrl+C -> finishes current frame, prints session summary.
KillSignal=SIGINT
TimeoutStopSec=60

# Don't loop hard if something is fundamentally wrong (camera busy, SSD missing, etc.).
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
